# scripts/hermes/systemd/awiki-pi5-reboot.service
[Unit]
Description=Scheduled Pi5 reboot (hang prevention — clears leaks, restarts all Umbrel apps)

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl reboot
